# git.update/git.update.cron.d
## /etc/cron.d/ EXAMPLE
## TODO SET PATH TO git.update.sh
*/5 * * * * root /path/git.update.sh >/dev/null 2>&1
